# deploy/secondlens.service
[Unit]
Description=secondlens bucket watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/secondlens
Environment=PYTHONUNBUFFERED=1
Environment=INTELLENS_CONFIG=/opt/secondlens/config/intellens/prod.yaml
ExecStart=/opt/secondlens/.venv/bin/python main.py --config config.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
