# deploy/secondlens.service
[Unit]
Description=secondlens bucket watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/secondlens
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 main.py --config config.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
